# system-files/.local/share/system-files/services/etc/runit/sv/windscribe/run
#!/bin/sh
sv check NetworkManager > /dev/null || exit 1
windscribe start
exec pwait -f 'windscribe start'
