# system-files/.local/share/system-files/hblock/etc/runit/sv/hblock/run
#!/bin/sh
while true; do
    sleep 1m
    /usr/local/bin/hblock-aria2 -l
done
