# services/.local/runit/sv/xss-lock/run
#!/bin/sh
exec xss-lock --transfer-sleep-lock -- sh -c '
    dunstctl set-paused true
    slock
    dunstctl set-paused false
'
