# services/.local/runit/sv/xss-lock/run
#!/bin/sh
( ( sleep 0.2; loginctl lock-session ) & )
exec xss-lock -- sh -c 'dunstctl set-paused true; slock; dunstctl set-paused false'
