# sysfiles/.local/share/sysfiles/hblock/hblock.cron
#!/bin/sh
# /etc/cron.daily/hblock
exec /usr/local/bin/hblock-aria2
